# config/backup/local_systemd/local-backup.service
[Unit]
Description=Rsync backup of home directory to ThinkPad Server

[Service]
Type=oneshot
ExecStart=/home/matanshalti/run_backup.sh
